# src/samples/flow/systemd/unix-socket-reader.service.in
[Unit]
Description=Test SystemCallFilter for unix socket (reader)

[Service]
ExecStart=@prefix@/bin/sol-fbp-runner @abs_top_builddir@/src/samples/flow/unix-socket/string-reader.fbp
